# net-misc/alfred-announce/files/alfred-announce.cron
# Global variables
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
HOME=/

* * * * * root /opt/alfred-announce/announce.sh -u /var/run/alfred.sock
